# supervisord.conf
[supervisord]
nodaemon=true

[program:jawanndenn]
command=jawanndenn --port 80 --host 0.0.0.0
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr